# sage/gateway/com.web4.mcnugget.probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.web4.mcnugget.probe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dennispalatov/repos/HRM/sage/scripts/mcnugget_probe_and_push.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/dennispalatov</string>
    </dict>

    <!-- Every 6 hours -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/Users/dennispalatov/Library/Logs/sage/probe.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dennispalatov/Library/Logs/sage/probe_error.log</string>

    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
