# sage/gateway/com.web4.mcnugget.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.web4.mcnugget.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/dennispalatov/repos/HRM/sage/scripts/mcnugget_sync.sh</string>
    </array>

    <!-- Every hour -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/dennispalatov/Library/Logs/sage/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dennispalatov/Library/Logs/sage/sync_error.log</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
